# ops/alienware-fan-cap.service
[Unit]
Description=Alienware persistent fan boost cap
After=multi-user.target
ConditionPathExists=/sys/class/hwmon/hwmon6/fan1_boost
ConditionPathExists=/sys/class/hwmon/hwmon6/fan2_boost

[Service]
Type=simple
ExecStart=/usr/local/sbin/alienware-fan-cap.sh 5000 5000 200 1
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
